# init/openrc_conf
#!/sbin/openrc-run

name="MineOS Minecraft Server Management Interface"

depend() {
        need net localmount
        after firewall
}

start() {
        start-stop-daemon --start \
                --chdir /usr/games/minecraft \
                --pidfile /var/run/mineos.pid \
                --exec "/usr/local/bin/ts-node" \
                -- ./src/service.ts start
}

stop() {
        start-stop-daemon --stop \
                --chdir /usr/games/minecraft \
                --pidfile /var/run/mineos.pid \
                --exec "/usr/local/bin/ts-node" \
                -- ./src/service.ts stop
}
